# pltd/profile.d/pltd
[X] jobname=set_kernel_ipv4_forward,template=linux_common/set_kernel_value|net.ipv4.ip_forward|0
[X] jobname=set_kernel_ipv4_default_send_redirects,template=linux_common/set_kernel_value|net.ipv4.conf.default.send_redirects|0
[X] jobname=set_kernel_ipv4_all_send_redirects,template=linux_common/set_kernel_value|net.ipv4.conf.all.send_redirects|0
[X] jobname=set_kernel_ipv4_all_source_route,template=linux_common/set_kernel_value|net.ipv4.conf.all.accept_source_route|0
[X] jobname=set_kernel_ipv4_default_source_route,template=linux_common/set_kernel_value|net.ipv4.conf.default.accept_source_route|0
[X] jobname=set_kernel_ipv4_all_accept_redirects,template=linux_common/set_kernel_value|net.ipv4.conf.all.accept_redirects|0
[X] jobname=set_kernel_ipv4_default_accept_redirects,template=linux_common/set_kernel_value|net.ipv4.conf.default.accept_redirects|0
[X] jobname=set_kernel_ipv4_all_secure_redirects,template=linux_common/set_kernel_value|net.ipv4.conf.all.secure_redirects|0
[X] jobname=set_kernel_ipv4_default_secure_redirects,template=linux_common/set_kernel_value|net.ipv4.conf.default.secure_redirects|0
[X] jobname=set_kernel_ipv4_all_log_martians,template=linux_common/set_kernel_value|net.ipv4.conf.all.log_martians|1
[X] jobname=set_kernel_ipv4_default_log_martians,template=linux_common/set_kernel_value|net.ipv4.conf.default.log_martians|1
[X] jobname=set_kernel_ipv4_icmp_echo_ignore_broadcasts,template=linux_common/set_kernel_value|net.ipv4.icmp_echo_ignore_broadcasts|1
[X] jobname=set_kernel_ipv4_icmp_ignore_bogus_error_responses,template=linux_common/set_kernel_value|net.ipv4.icmp_ignore_bogus_error_responses|1
[X] jobname=set_kernel_ipv4_all_rp_filter,template=linux_common/set_kernel_value|net.ipv4.conf.all.rp_filter|1
[X] jobname=set_kernel_ipv4_default_rp_filter,template=linux_common/set_kernel_value|net.ipv4.conf.default.rp_filter|1
[X] jobname=set_kernel_ipv4_tcp_syncookies,template=linux_common/set_kernel_value|net.ipv4.tcp_syncookies|1
[X] jobname=set_kernel_ipv6_all_accept_ra,template=linux_common/set_kernel_value|net.ipv6.conf.all.accept_ra|0
[X] jobname=set_kernel_ipv6_default_accept_ra,template=linux_common/set_kernel_value|net.ipv6.conf.default.accept_ra|0
[X] jobname=set_kernel_ipv6_all_accept_redirects,template=linux_common/set_kernel_value|net.ipv6.conf.all.accept_redirects|0
[X] jobname=set_kernel_ipv6_default_accept_redirects,template=linux_common/set_kernel_value|net.ipv6.conf.default.accept_redirects|0
[X] jobname=set_service_off_chargen_stream,template=pontis_rhel6/set_service6_disable|chargen-stream
[X] jobname=set_service_off_daytime_dgram,template=pontis_rhel6/set_service6_disable|daytime-dgram
[X] jobname=set_service_off_chargen_dgram,template=pontis_rhel6/set_service6_disable|chargen-dgram
[X] jobname=set_service_off_daytime_stream,template=pontis_rhel6/set_service6_disable|daytime-stream
[X] jobname=set_service_off_discard_stream,template=pontis_rhel6/set_service6_disable|discard-stream
[X] jobname=set_service_off_discard_dgram,template=pontis_rhel6/set_service6_disable|discard-dgram
[X] jobname=set_service_off_echo_dgram,template=pontis_rhel6/set_service6_disable|echo-dgram
[X] jobname=set_service_off_echo_stream,template=pontis_rhel6/set_service6_disable|echo-stream
[X] jobname=set_service_off_time_stream,template=pontis_rhel6/set_service6_disable|time-stream
[X] jobname=set_service_off_time_dgram,template=pontis_rhel6/set_service6_disable|time-dgram
[X] jobname=set_service_off_tftp,template=pontis_rhel6/set_service6_disable|tftp
[X] jobname=set_service_off_xinetd,template=pontis_rhel7/set_service7_disable|xinetd
[X] jobname=set_service_off_auditd,template=pontis_rhel7/set_service7_disable|auditd
