# LHT-2.2/profile.d/profile
[X] jobname=mybew,template=linux_common/create_dir|/home/toto
[X] jobname=myuse,template=linux_common/create_user|toto|10000|10000
