# LHT-2.2/profile.d/update_patch
[X] jobname=check_patch_update,template=linux_common/check_patch_update
