# pontis/profile.d/pontis_rhel6
[X] jobname=create_group_mapr,template=linux_common/create_group|mapr|1000
[X] jobname=create_user_mapr,template=linux_common/create_user|mapr|1000|1000
[X] jobname=create_user_pontis,template=linux_common/create_user|pontis|10025|1000
[X] jobname=create_user_PontisSFTP,template=linux_common/create_user|PontisSFTP||
[X] jobname=delete_user_lp,template=linux_common/delete_user|lp
[X] jobname=delete_group_lp,template=linux_common/delete_group|lp
[X] jobname=delete_user_uucp,template=linux_common/delete_user|uucp
[X] jobname=delete_group_uucp,template=linux_common/delete_group|uucp
[X] jobname=delete_user_games,template=linux_common/delete_user|games
[X] jobname=delete_group_games,template=linux_common/delete_group|games
[X] jobname=delete_user_gopher,template=linux_common/delete_user|gopher
[X] jobname=delete_group_gopher,template=linux_common/delete_group|gopher
[X] jobname=delete_user_news,template=linux_common/delete_user|news
[X] jobname=delete_user_mailnull,template=linux_common/delete_user|mailnull
[X] jobname=delete_user_adm,template=linux_common/delete_user|adm
[X] jobname=create_group_sftponly,template=linux_common/create_group|sftponly
[X] jobname=delete_group_adm,template=linux_common/delete_group|adm
[X] jobname=delete_group_dip,template=linux_common/delete_group|dip
[X] jobname=add_PontisSFTP_into_sftponly,template=linux_common/add_user_to_group|PontisSFTP|sftponly
[X] jobname=set_shell_user_sync,template=linux_common/set_user_shell|sync|/sbin/nologin
[X] jobname=set_shell_user_mail,template=linux_common/set_user_shell|mail|/sbin/nologin
[X] jobname=lock_password_user_sync,template=linux_common/lock_user_password|sync
[X] jobname=lock_password_user_mail,template=linux_common/lock_user_password|mail
[X] jobname=create_file_sftp_config,template=linux_common/create_file|/etc/ssh/sftp_config
[X] jobname=create_dir_sftp,template=linux_common/create_dir|/var/sftp
[X] jobname=create_dir_sftp_uploads,template=linux_common/create_dir|/var/sftp/uploads
[X] jobname=set_dir_sftp_owner,template=linux_common/set_dir_owner|/var/sftp|root
[X] jobname=set_dir_sftp_group,template=linux_common/set_dir_group|/var/sftp|root
[X] jobname=set_dir_sftp_perm,template=linux_common/set_dir_perm|/var/sftp|755
[X] jobname=set_dir_uploads_owner,template=linux_common/set_dir_owner|/var/sftp/uploads|PontisSFTP
[X] jobname=set_dir_uploads_group,template=linux_common/set_dir_group|/var/sftp/uploads|PontisSFTP
[X] jobname=set_file_securetty_owner,template=linux_common/set_file_owner|/etc/securetty|root
[X] jobname=set_file_securetty_group,template=linux_common/set_file_group|/etc/securetty|root
[X] jobname=set_dir_tmp_perm,template=linux_common/set_dir_perm|/tmp|1777
[X] jobname=set_dir_var_tmp_perm,template=linux_common/set_dir_perm|/var/tmp|1777
[X] jobname=set_file_wtmp_owner,template=linux_common/set_file_owner|/var/log/wtmp|root
[X] jobname=set_file_wtmp_group,template=linux_common/set_file_group|/var/log/wtmp|root
[X] jobname=set_file_wtmp_perm,template=linux_common/set_file_perm|/var/log/wtmp|600
[X] jobname=set_file_grub_conf_perm,template=linux_common/set_file_perm|/etc/grub.conf|600
[X] jobname=set_banner_issue_content,template=pontis_custom/update_banner_content|/etc/issue|./misc/banner
[X] jobname=set_banner_issue_net_content,template=pontis_custom/update_banner_content|/etc/issue.net|./misc/banner
[X] jobname=set_banner_motd_content,template=pontis_custom/update_banner_content|/etc/motd|./misc/banner
[X] jobname=set_root_umask_0027,template=pontis_custom/update_umask_bashrc
[X] jobname=set_user_shell_timeout,template=linux_common/set_user_shell_timeout|900
[X] jobname=set_ssh_loglevel,template=pontis_custom/set_ssh_loglevel
[X] jobname=set_ssh_protocol,template=pontis_custom/set_ssh_protocol
[X] jobname=set_ssh_permit_root_login,template=pontis_custom/set_ssh_permit_root_login
[X] jobname=set_ssh_max_auth_tries,template=pontis_custom/set_ssh_max_auth_tries
[X] jobname=set_ssh_empty_password,template=pontis_custom/set_ssh_empty_password
[X] jobname=set_ssh_alive_count,template=pontis_custom/set_ssh_alive_count
[X] jobname=set_ssh_alive_interval,template=pontis_custom/set_ssh_alive_interval
[X] jobname=set_ssh_port,template=pontis_custom/set_ssh_port
[X] jobname=set_ssh_ciphers,template=pontis_custom/set_ssh_ciphers
[X] jobname=set_ssh_rhosts,template=pontis_custom/set_ssh_rhosts
[X] jobname=set_ssh_host_based_auth,template=pontis_custom/set_ssh_host_based_auth
[X] jobname=set_ssh_banner,template=pontis_custom/set_ssh_banner
[X] jobname=set_ssh_permit_user_env,template=pontis_custom/set_ssh_permit_user_env
[X] jobname=set_ssh_strict_mode,template=pontis_custom/set_ssh_strict_mode
[X] jobname=set_ssh_privilege_separation,template=pontis_custom/set_ssh_privilege_separation
[X] jobname=set_ssh_allow_users,template=pontis_custom/set_ssh_allow_users
[X] jobname=disable_selinux,template=linux_common/disable_selinux
[X] jobname=set_auth_passwd_history,template=pontis_custom/set_auth_passwd_history
[X] jobname=set_passwd_min_days,template=pontis_custom/set_passwd_min_days
[X] jobname=set_auth_enforce_root,template=pontis_rhel6/set_auth_enforce_root
[X] jobname=set_auth_passwd_min_len,template=pontis_rhel6/set_auth_passwd_min_len
[X] jobname=set_auth_passwd_min_class,template=pontis_rhel6/set_auth_passwd_min_class
[X] jobname=set_ctrl_alt_delete_override,template=pontis_rhel6/set_ctrl_alt_delete_override
[X] jobname=set_init_runlevel_3,template=pontis_rhel6/set_init_runlevel_3
[X] jobname=set_sudo_syslog_auth,template=pontis_custom/set_sudo_syslog_auth
[X] jobname=set_sudo_pontis,template=pontis_custom/set_sudo_pontis
[X] jobname=set_sudo_mapr,template=pontis_custom/set_sudo_mapr
[X] jobname=set_rsh_off_securetty,template=pontis_custom/set_rsh_off_securetty
[X] jobname=set_pts_off_securetty,template=pontis_custom/set_pts_off_securetty
[X] jobname=add_logrotate_cron,template=pontis_custom/add_logrotate_cron
[X] jobname=add_logrotate_messages,template=pontis_custom/add_logrotate_messages
[X] jobname=add_logrotate_secure,template=pontis_custom/add_logrotate_secure
[X] jobname=add_logrotate_spooler,template=pontis_custom/add_logrotate_spooler
[X] jobname=add_logrotate_maillog,template=pontis_custom/add_logrotate_maillog
[X] jobname=add_logrotate_bootlog,template=pontis_custom/add_logrotate_bootlog
[X] jobname=disable_service_gpm,template=pontis_rhel6/set_service6_disable|gpm
[X] jobname=disable_service_kudzu,template=pontis_rhel6/set_service6_disable|kudzu
[X] jobname=disable_service_apmd,template=pontis_rhel6/set_service6_disable|apmd
[X] jobname=disable_service_rwhod,template=pontis_rhel6/set_service6_disable|rwhod
[X] jobname=disable_service_isdn,template=pontis_rhel6/set_service6_disable|isdn
[X] jobname=disable_service_rstatd,template=pontis_rhel6/set_service6_disable|rstatd
[X] jobname=disable_service_rusersd,template=pontis_rhel6/set_service6_disable|rusersd
[X] jobname=disable_service_vncserver,template=pontis_rhel6/set_service6_disable|vncserver
[X] jobname=disable_service_arpwatch,template=pontis_rhel6/set_service6_disable|arpwatch
[X] jobname=disable_service_vsftpd,template=pontis_rhel6/set_service6_disable|vsftpd
[X] jobname=disable_service_xfs,template=pontis_rhel6/set_service6_disable|xfs
[X] jobname=disable_service_nfslock,template=pontis_rhel6/set_service6_disable|nfslock
[X] jobname=disable_service_squid,template=pontis_rhel6/set_service6_disable|squid
[X] jobname=disable_service_atd,template=pontis_rhel6/set_service6_disable|atd
[X] jobname=disable_service_rhnsd,template=pontis_rhel6/set_service6_disable|rhnsd
[X] jobname=disable_service_telnet,template=pontis_rhel6/set_service6_disable|telnet
[X] jobname=disable_service_amanda,template=pontis_rhel6/set_service6_disable|amanda
[X] jobname=disable_service_daytime_dgram,template=pontis_rhel6/set_service6_disable|daytime-dgram
[X] jobname=disable_service_cvs,template=pontis_rhel6/set_service6_disable|cvs
[X] jobname=disable_service_krb5_telnet,template=pontis_rhel6/set_service6_disable|krb5-telnet
[X] jobname=disable_service_rexec,template=pontis_rhel6/set_service6_disable|rexec
[X] jobname=disable_service_tftp,template=pontis_rhel6/set_service6_disable|tftp
[X] jobname=disable_service_amandaidx,template=pontis_rhel6/set_service6_disable|amandaidx
[X] jobname=disable_service_daytime_stream,template=pontis_rhel6/set_service6_disable|daytime-stream
[X] jobname=disable_service_gssftp,template=pontis_rhel6/set_service6_disable|gssftp
[X] jobname=disable_service_rlogin,template=pontis_rhel6/set_service6_disable|rlogin
[X] jobname=disable_service_talk,template=pontis_rhel6/set_service6_disable|talk
[X] jobname=disable_service_time_dgram,template=pontis_rhel6/set_service6_disable|time-dgram
[X] jobname=disable_service_ntalk,template=pontis_rhel6/set_service6_disable|ntalk
[X] jobname=disable_service_amidxtape,template=pontis_rhel6/set_service6_disable|amidxtape
[X] jobname=disable_service_rsh,template=pontis_rhel6/set_service6_disable|rsh
[X] jobname=disable_service_eklogin,template=pontis_rhel6/set_service6_disable|eklogin
[X] jobname=disable_service_auth,template=pontis_rhel6/set_service6_disable|auth
[X] jobname=disable_service_echo-stream,template=pontis_rhel6/set_service6_disable|echo-stream
[X] jobname=disable_service_ktalk,template=pontis_rhel6/set_service6_disable|ktalk
[X] jobname=disable_service_discard-dgram,template=pontis_rhel6/set_service6_disable|discard-dgram
[X] jobname=disable_service_klogin,template=pontis_rhel6/set_service6_disable|klogin
[X] jobname=disable_service_time_stream,template=pontis_rhel6/set_service6_disable|time-stream
[X] jobname=disable_service_discard_stream,template=pontis_rhel6/set_service6_disable|discard-stream
[X] jobname=set_sftp_server_force_command,template=pontis_custom/set_sftp_server_force_command
[X] jobname=set_sftp_server_allow_agent_forwarding,template=pontis_custom/set_sftp_server_allow_agent_forwarding
[X] jobname=set_sftp_server_allow_tcp_forwarding,template=pontis_custom/set_sftp_server_allow_tcp_forwarding
[X] jobname=set_sftp_server_password_auth,template=pontis_custom/set_sftp_server_password_auth
[X] jobname=set_sftp_server_chal_response_auth,template=pontis_custom/set_sftp_server_chal_response_auth
[X] jobname=set_sftp_server_loglevel,template=pontis_custom/set_sftp_server_loglevel
[X] jobname=set_sftp_server_permit_tunnel,template=pontis_custom/set_sftp_server_permit_tunnel
[X] jobname=set_sftp_subsystem_server,template=pontis_custom/set_sftp_subsystem_server
[X] jobname=set_sftp_match_group,template=pontis_custom/set_sftp_match_group
[X] jobname=set_sftp_server_chroot,template=pontis_custom/set_sftp_server_chroot
